# tools-internal/kdc-api-server/docker/supervisord.conf
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=102
minprocs=200
user=root

[program:kdc]
command=%(ENV_KDC_BIN) -r LOCAL -n
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0

[program:webapi]
command=/kdc/api-server
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
